# src/hermes_companion/systemd/hermes-companion.service
[Unit]
Description=Hermes Companion Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 -m companion_cli serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
